# os/overlay/etc/systemd/system/krdos-maintenance.service
[Unit]
Description=KrdOS Weekly Maintenance
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/maintenance.sh
StandardOutput=journal
StandardError=journal
